# photobooth-kiosk.service
[Unit]
Description=Photobooth Kiosk Mode
Documentation=https://github.com/Glebastik/photobooth_dd
After=multi-user.target graphical-session.target
Wants=multi-user.target
Conflicts=lightdm.service gdm.service

[Service]
Type=simple
User=ddself
Group=ddself
WorkingDirectory=/opt/photobooth
ExecStart=/usr/local/bin/start-kiosk.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Увеличенный timeout для запуска X сервера
TimeoutStartSec=60
TimeoutStopSec=30

# Переменные окружения
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XAUTHORITY=/home/ddself/.Xauthority
Environment=HOME=/home/ddself

# Права доступа к устройствам
SupplementaryGroups=video audio input dialout tty

# Безопасность
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# Ограничения ресурсов
MemoryLimit=2G
CPUQuota=80%

# Создание директории логов в домашней папке (безопасно)
ExecStartPre=/bin/mkdir -p /home/ddself/logs

[Install]
WantedBy=multi-user.target
